# scripts/com.safar.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.safar.daily</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/shaheersaud/Safar/.venv/bin/python</string>
		<string>/Users/shaheersaud/Safar/safar_daily_pipeline.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/shaheersaud/Safar</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/shaheersaud/Safar/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>6</integer>
		<key>Minute</key>
		<integer>15</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/shaheersaud/Library/Logs/safar_pipeline.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/shaheersaud/Library/Logs/safar_pipeline.err</string>
</dict>
</plist>
